# systemd/user/eve-av-bridge-broadcast.service
[Unit]
Description=eve AV bridge (stream-broadcast: freeq call → stream.place RTMP)
Documentation=file://@ROOT@/systemd/README.md
After=network-online.target eve-av-bridge.service
PartOf=eve.target

[Service]
Type=simple
WorkingDirectory=@ROOT@
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=STREAM_BROADCAST_AV_BRIDGE_BIND=127.0.0.1:8793
Environment=AV_BRIDGE_BIND=127.0.0.1:8793
Environment=AV_PLANE_ROLE=broadcast
Environment=RADIO_TITLE_HOOK=
EnvironmentFile=-%h/.config/eve/runtime.env
EnvironmentFile=-%h/.config/eve/config.env
ExecStart=@ROOT@/scripts/run-av-bridge-broadcast.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=eve.target
